# deploy/ubuntu24/systemd/jijin-api.service
[Unit]
Description=Jijin FastAPI Backend
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/jijin/backend
EnvironmentFile=/opt/jijin/backend/.env
ExecStart=/opt/jijin/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
